# rootfs/etc/s6-overlay/s6-rc.d/api2sbs/run
#!/bin/sh
exec /etc/s6-overlay/scripts/api2sbs
